# s6-rc.d/sshd_init/up
#!/command/execlineb -P
bash /setup_sshd.sh